# backup.timer
[Unit]
Description=periodically backup home and storage mounts

[Timer]
OnCalendar=hourly
Persistent=True

[Install]
WantedBy=timers.target
