# hubble/conf/etc/systemd/system/backup.timer
[Unit]
Description=periodic backup of home folder

[Timer]
OnCalendar=hourly
Persistent=True

[Install]
WantedBy=timers.target
