# galileo/conf/etc/systemd/system/getty@tty11.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/dmesg --decode --human --time-format ctime
StandardInput=tty
StandardOutput=tty
TimeoutStopSec=0.1
