# galileo/conf/etc/systemd/system/getty@tty12.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/top
StandardInput=tty
StandardOutput=tty
